# daemons/expire_token.service
[Service]
Type=oneshot
WorkingDirectory=/var/www/equipe500/projet6
ExecStart=php daemons/expire_token.php
